# src/etc/cron.d/alternc-borgbackup
0 2 * * *       alterncpanel    /usr/lib/alternc/borgbackup.sh
